# services/receive_operations/receive_operations.socket
[Unit]
Description=Receive Operations Socket
PartOf=receive_operations.service

[Socket]
ListenStream=0.0.0.0:8888

[Install]
WantedBy=sockets.target